# .config/systemd/user/bible.service
[Unit]
Description=Generate random Bible verse

[Service]
Type=oneshot
ExecStart=/home/fera/.local/bin/bible_daemon
