# .config/systemd/user/bible.timer
[Unit]
Description=Refresh Bible verse every 5 seconds

[Timer]
OnBootSec=2
OnUnitActiveSec=5
AccuracySec=1

[Install]
WantedBy=default.target
